# .config/systemd/user/random-wallpaper.service
[Unit]
Description=Run random-wallpaper.sh

[Service]
ExecStart=%h/.config/scripts/random-wallpaper.sh
